# my-deb-template/DEBIAN/postinst
sudo systemctl enable robot
sudo systemctl start robot
